# deploy/idealimage-celerybeat.service
[Unit]
Description=IdealImage Celery Beat Scheduler
After=network.target redis.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/idealimage
EnvironmentFile=/var/www/idealimage/.env
ExecStart=/var/www/idealimage/venv/bin/celery \
    -A IdealImage_PDJ \
    beat \
    -l info \
    --scheduler django_celery_beat.schedulers:DatabaseScheduler \
    --logfile=/var/log/idealimage/celery-beat.log
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
